# xrootd/xrootd-monitoring-shoveler/config/default/etc/supervisord.d/10-xrootd-monitoring-shoveler.conf
[program:xrootd-monitoring-shoveler]
command=/usr/bin/xrootd-monitoring-shoveler
user=xrootd
group=xrootd
autorestart=true

