# xrootd/xcache-mon/config/default/etc/supervisord.d/10-xrootd.conf
[program:xcache-mon]
command=python3 /opt/t2_maitenance/src/python/t2Mon/XRootD/xrootdcache.py
user=xrootd
autorestart=true
