# xrootd/xcache-server/config/default/etc/supervisord.d/10-xrootd.conf
[program:startup]
command=/usr/local/sbin/startup.sh
priority=-1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=false

[program:server-cmsd]
command=bash -c "sleep 5 && cmsd -c /etc/xrootd/xrootd.cfg -l /var/log/xrootd/cmsd.log"
user=xrootd
autorestart=true

[program:server-xrootd]
command=bash -c "sleep 5 && xrootd -c /etc/xrootd/xrootd.cfg -l /var/log/xrootd/xrootd.log"
user=xrootd
autorestart=true

[program:xrootd-mon]
command=python3 /opt/xroot-mon.py
autorestart=true

[program:gunicorn]
command=bash -c "cd /opt; gunicorn --workers=2  http-server:app -b 0.0.0.0:8000"
autorestart=true

#[program:xcache-rucio-reporter]
#command=bash -c "sleep 5 && /opt/xcache/run.sh"
#autorestart=true
