# xrootd/xrootd-clustered-server/config/default/etc/supervisord.d/10-xrootd.conf
[program:clustered-cmsd]
command=cmsd -c /etc/xrootd/xrootd-clustered.cfg -l /var/log/xrootd/clustered/cmsd.log
user=xrootd
autorestart=true

[program:clustered-xrootd]
command=xrootd -c /etc/xrootd/xrootd-clustered.cfg -l /var/log/xrootd/clustered/xrootd.log
user=xrootd
autorestart=true

[program:xrootd-mon]
command=python3 /opt/t2_maitenance/src/python/t2Mon/XRootD/xrootd.py
user=xrootd
autorestart=true
