# xrootd/xrootd-stageout-server/config/default/etc/supervisord.d/10-xrootd.conf
[program:clustered-cmsd]
command=bash -c "sleep 15 && cmsd -c /etc/xrootd/xrootd-clustered.cfg -l /var/log/xrootd/clustered/cmsd.log"
user=xrootd
autorestart=true

[program:clustered-xrootd]
# Note "-R xrootd" here instructs xrootd to drop privileges to the xrootd Unix user.
command=bash -c "sleep 15 && xrootd -c /etc/xrootd/xrootd-clustered.cfg -l /var/log/xrootd/clustered/xrootd.log -R xrootd"
autorestart=true
# These provide xrootd with the ability to switch UIDs/GIDs for reading/writing files.
CapabilityBoundingSet=CAP_SETUID CAP_SETGID CAP_DAC_OVERRIDE
Capabilities=CAP_SETGID+p CAP_SETUID+p

[program:xrootd-mon]
command=python3 /opt/xroot-mon.py
autorestart=true

[program:gunicorn]
command=bash -c "cd /opt; gunicorn --workers=2  http-server:app -b 0.0.0.0:8000"
autorestart=true
