# wlcg-site-mon/config/10-wlcg-mon.conf
[program:wlcg-site-snmp-mon]
command=/bin/bash -c "cd /opt/ && python3 WLCG-site-snmp-http.py"
autostart=true
autorestart=true
