# xrootd/xrootd-sdn/redir/config/default/etc/supervisord.d/10-xrootd.conf
[program:clustered-cmsd]
command=cmsd -I v6 -c /etc/xrootd/xrootd-clustered.cfg -l /var/log/xrootd/clustered/cmsd.log
user=xrootd
autorestart=true

[program:clustered-xrootd]
command=xrootd -I v6 -c /etc/xrootd/xrootd-clustered.cfg -l /var/log/xrootd/clustered/xrootd.log
user=xrootd
autorestart=true
